# systemd/wintarch-first-boot.service
[Unit]
Description=Wintarch First Boot Configuration
After=local-fs.target dbus.service
# Keep first-boot configuration ahead of SDDM and the generic display manager alias.
Before=sddm.service display-manager.service
ConditionPathExists=!/var/lib/wintarch/first-boot-done

[Service]
Type=oneshot
ExecStart=/opt/wintarch/install/first-boot.sh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
